# rc.local
#!/bin/sh -e
#
#rc.local
#

bash /root/MO_server.sh

exit 0
